# config/systemd/user/tesseract-web-tunnel.service
[Unit]
Description=Tesseract Web Cloudflare Tunnel
Wants=network-online.target
Requires=tesseract-web-production.service
After=network-online.target tesseract-web-production.service

[Service]
Type=simple
Environment=HOME=/home/bot
ExecStart=/home/bot/.local/bin/cloudflared --no-autoupdate tunnel run --token-file /home/bot/.config/tesseract/tesseract-web-tunnel.token
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGTERM
UMask=0077

[Install]
WantedBy=default.target
