# deploy/tesseract-deploy.service
# Systemd oneshot service that runs bin/deploy
# Triggered by tesseract-deploy.timer
# Install: sudo cp deploy/tesseract-deploy.service /etc/systemd/system/
#          sudo systemctl daemon-reload

[Unit]
Description=Tesseract auto-deploy
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/user/tesseract
ExecStart=/home/user/tesseract/bin/deploy
EnvironmentFile=/home/user/tesseract/.env
