# deploy/tesseract-deploy.timer
# Systemd timer that polls for new commits and deploys
# Install: sudo cp deploy/tesseract-deploy.timer /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now tesseract-deploy.timer

[Unit]
Description=Poll for Tesseract deploys every minute

[Timer]
OnBootSec=30
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
